# scripts/crontab
# mysql backup
0 0 * * * /vagrant/scripts/mysqlbak.sh
